# nlxtp_guest_modules/linux/debian/nlxtp-guest-modules-0.1/debian/nlxtp-guest-modules/DEBIAN/postrm
#!/bin/bash
update-rc.d -f nlXTP_daemon remove

