# docker/build/tomato_service/inittab
::sysinit:/init.sh >/dev/stdout 2>/dev/stderr
::respawn:/code/service/server.py >/dev/stdout 2>/dev/stderr

