# dist/nlxtp_guest_modules/custom_files/debian/postinst
#!/bin/bash
update-rc.d nlXTP_daemon defaults

