# nlxtp_guest_modules/linux/debian/nlxtp-guest-modules-0.1/debian/postinst
#!/bin/bash
update-rc.d nlXTP_daemon defaults
echo 33 > /success.tmp
