# docker/build/tomato_backend_accounting/inittab
::sysinit:/init.sh >/dev/stdout 2>/dev/stderr
::respawn:/code/server.sh >/dev/stdout 2>/dev/stderr
