# docker/build/tomato-backend/inittab
::sysinit:/init.sh >/dev/stdout 2>/dev/stderr
::respawn:/code/server.py >/dev/stdout 2>/dev/stderr

